# ads-registry.service
[Unit]
Description=ADS Container Registry
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ads-registry
ExecStart=/opt/ads-registry/ads-registry serve --config /opt/ads-registry/config.json
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ads-registry

# Security
NoNewPrivileges=true
PrivateTmp=true

# Limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
